# hosting/com.meandlia.geo-backfill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.meandlia.geo-backfill</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/Foongbear/melia-geo-backfill.py</string>
  </array>
  <key>StartInterval</key><integer>3600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/Foongbear/Library/Logs/melia-geo-backfill.log</string>
  <key>StandardErrorPath</key><string>/Users/Foongbear/Library/Logs/melia-geo-backfill.log</string>
</dict>
</plist>
